# lib/systemd/system/concat.service
[Unit]
Description=Concat
[Service]
Type=oneshot
ExecStart=/var/www/netaggregator.ru/src/App/Service/Stream/concat.sh
PrivateTmp=true
