# etc/cron.d/renew
0 */12 * * * /var/www/netaggregator.ru/src/App/Service/Stream/concat.sh
